# deploy/pantri.service
[Unit]
Description=Pantri API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/pantri
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin

[Install]
WantedBy=multi-user.target
